# etc/profile.d/set-java-home.sh
JAVA_HOME=$(readlink -f $(which javac) | sed 's|/bin/javac||')
